# scripts/openclaw-dashboard.service
[Unit]
Description=OpenClaw Dashboard
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/openclaw-dashboard
Environment="PATH=/opt/openclaw-dashboard/venv/bin"
Environment="PYTHONPATH=/opt/openclaw-dashboard"
Environment="OPENCLAW_PATH=/root/.openclaw"
Environment="FLASK_APP=app.py"
ExecStart=/opt/openclaw-dashboard/venv/bin/python -m flask run --host=0.0.0.0 --port=8088
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
